# deploy/systemd/cloudflared-tunnel.service
[Unit]
Description=Cloudflared Quick Tunnel for HomeVault
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/xe3
Environment=PYTHONPATH=%h/xe3
ExecStart=%h/xe3/venv/bin/python %h/xe3/scripts/cloudflared_tunnel.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
