# deploy/systemd/cloudflared-watchdog.service
[Unit]
Description=Cloudflared Tunnel Watchdog for HomeVault
After=network-online.target cloudflared-tunnel.service
Wants=network-online.target cloudflared-tunnel.service

[Service]
Type=simple
WorkingDirectory=%h/xe3
Environment=PYTHONPATH=%h/xe3
EnvironmentFile=%h/xe3/.env
ExecStart=%h/xe3/venv/bin/python %h/xe3/scripts/tunnel_watchdog.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
