# deploy/systemd/xe3-web.service
[Unit]
Description=XE3 Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/xe3
Environment=PYTHONPATH=%h/xe3
EnvironmentFile=%h/xe3/.env
ExecStart=%h/xe3/venv/bin/python %h/xe3/apps/web/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
